# base/mkosi.skeleton/etc/systemd/system/network-setup.service
[Unit]
Description=Basic Network Setup
DefaultDependencies=no
Before=network.target
Wants=network.target

[Service]
Type=oneshot
ExecStart=ip link set lo up
ExecStart=ip link set eth0 up
ExecStart=chattr +i /etc/resolv.conf
ExecStart=/usr/sbin/udhcpc -i eth0 -n
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target