# services/systemd/rbuilder.service
[Unit]
Description=RBuilder Bidding Service
After=network.target network-setup.service persistent-mount.service
Requires=network-setup.service persistent-mount.service

[Service]
Type=exec
User=rbuilder
Group=eth
ExecStartPre=+/usr/bin/rbuilder-init
ExecStart=/usr/bin/bwrap \
    --ro-bind /usr /usr \
    --ro-bind /lib /lib \
    --ro-bind /lib64 /lib64 \
    --ro-bind /bin /bin \
    --ro-bind /sbin /sbin \
    --ro-bind /etc/bidding.toml /config.toml \
    --bind /var/run/rbuilder /var/run/rbuilder \
    --proc /proc \
    --dev /dev \
    --clearenv \
    --unshare-pid \
    bidding-service /config.toml
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=minimal.target