# mkosi.skeleton/etc/inittab
id:2:initdefault:
si::sysinit:/etc/init.d/rcS
s0:2345:respawn:/bin/sh